# etc/xdg/autostart/pika-login-sound.desktop
[Desktop Entry]
Name=Play Login Sound on startup
Exec=/etc/pika_login.sh
Type=Application
Terminal=false
